# deploy/prayerhub.service
[Unit]
Description=PrayerHub Mini Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/prayerhub
ExecStart=/opt/prayerhub/.venv/bin/python -m prayerhub
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
